# systemd/cuebeam.service
[Unit]
Description=CueBeam Player
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/CueBeam
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/pi/CueBeam/.venv/bin/python /home/pi/CueBeam/app.py
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
